# systemd/fuse-dynamic-mount.service
[Service]
Type=oneshot
ExecStart=/bin/dynamic-mount all --notify

[Install]
WantedBy=user-applications.target
